# profile.d/development-with-javadocs
Javadocs required by downstream projects even when in development profile
